# systemd/sync-mail.timer
[Unit]
Description=Synchronize mailboxes regularly

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=default.target
