# systemd/sync-mail.service
[Unit]
Description=Synchronize mailboxes
Documentation=man:mbsync(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync --all
